# 99-usb-audio-update.rules
ACTION=="add", SUBSYSTEM=="sound", KERNEL=="card*", RUN+="/usr/local/bin/update_audio.sh"
